# rules/Dev_Sigma_win_scheduled_task_deletion.yml
rule  Dev_sigma_scheduled_task_deletion {
 meta:
    author = "David Strassegger"
    description = "Detects scheduled task deletion events. Scheduled tasks are likely to be deleted if not used for persistence. Malicious Software often creates tasks directly under the root node e.g. \\TASKNAME Author: David Strassegger."
    reference = ""
    version = "0.01"
    created = "2021/01/22"
    product = "windows"
    service = "security"
    mitre = "execution, privilege_escalation, t1053, car.2013-08-001, t1053.005"

  events:
$selection.metadata.product_event_type = "4699"

  condition:
    $selection
}
